# system/etc/udev/rules.d/99-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/meroving/.local/bin/scripts/connect_monitor"

